# oms-xdr/deploy/oms-netscan-full.timer
[Unit]
Description=OMS-XDR - scan complet hebdomadaire (dimanche 03:00)

[Timer]
OnCalendar=Sun *-*-* 03:00:00
Persistent=true
Unit=oms-netscan-full.service

[Install]
WantedBy=timers.target
