# oms-graph/deploy/oms-graph-respond.timer
[Unit]
Description=OMS-GRAPH - gradation des leurres declenches (toutes les 15 min, audit seul)

[Timer]
OnBootSec=14min
OnUnitActiveSec=15min
AccuracySec=2min
Unit=oms-graph-respond.service

[Install]
WantedBy=timers.target
