# oms-graph/deploy/oms-graph.service
[Unit]
Description=OMS-GRAPH - jumeau d'attaque / analyse d'exposition passive (chemins vers joyaux)
After=network-online.target opensearch.service graylog-server.service
Wants=network-online.target

[Service]
Type=oneshot
User=root
WorkingDirectory=/root/omnitech-siem-setup/oms-graph
ExecStart=/root/omnitech-siem-setup/oms-graph/.venv/bin/python -m oms_graph.run analyze --push --config /etc/oms-graph/config.yaml
Nice=15
