# oms-ml/deploy/oms-ml-fp.service
[Unit]
Description=OMS-ML - (re)entrainement du modele de reduction de faux positifs (supervise)
After=network-online.target opensearch.service
Wants=network-online.target

[Service]
Type=oneshot
User=root
WorkingDirectory=/root/omnitech-siem-setup/oms-ml
# S'auto-saute tant qu'il n'y a pas assez de labels analystes (cas qualifies VP/FP).
ExecStart=/root/omnitech-siem-setup/oms-ml/.venv/bin/python -m oms_ml.run fp --train --config /etc/oms-ml/config.yaml
Nice=10
